# .config/systemd/user/auto-commit-grimuar.service
[Unit]
Description=Auto Commit Grimuar
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/home/mia/Гримуар/Tools/Scripts/auto-commit.sh
